# crontab/crontab_settings
0 8 * * 7 rkhunter -c
0 9 * * 7 lynis audit system